# scheduler/clippertv-scheduler.service
[Unit]
Description=ClipperTV Monthly PDF Ingestion Scheduler
After=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/schedule-ingestion
Environment="PATH=/home/pi/.local/bin:/usr/local/bin:/usr/bin:/bin"
# Optional: Enable Healthchecks.io monitoring (uncomment and set your URL)
# Environment="HEALTHCHECK_URL=https://hc-ping.com/your-uuid-here"

# Run scheduler using uv
ExecStart=/home/pi/.local/bin/uv run python -m clippertv.scheduler.service

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=clippertv-scheduler

# Security hardening (optional)
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
